# crontab
*/5 * * * * cd /app && python3 manage.py delete_expired_reservations >> /app/cron_execution.log 2>&1
